# .config/systemd/user/cnf-fetch.service
[Unit]
Description=dotfiles updater

[Service]
Type=oneshot
ExecStart=%h/bin/cnf fetch
